# crontab
*/1 * * * * node /src/test-runner.js >> /var/log/e2e-tests.log 2>&1