# crontab
0 */6 * * * cd /src && /usr/local/bin/node --env-file=.env test-runner.js >> /var/log/e2e-test.log 2>&1
